# debian/postinst
#!/bin/sh

systemctl enable epipynet.service

systemctl restart epipynet.service
