# debian/prerm
#!/bin/sh

systemctl stop epipynet.service

systemctl disable epipynet.service
